# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/data/logs

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ============================================================
; Xvfb - Virtual display for Computer Use
; ============================================================
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -nolisten tcp
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/xvfb.log
stderr_logfile=/data/logs/xvfb.err
environment=HOME="/root"

; ============================================================
; Fluxbox - Window manager (required for browser automation)
; ============================================================
[program:fluxbox]
command=fluxbox -display :99
autostart=true
autorestart=true
priority=20
startretries=3
stdout_logfile=/data/logs/fluxbox.log
stderr_logfile=/data/logs/fluxbox.err
environment=DISPLAY=":99"

; ============================================================
; FastAPI Dashboard (installer/app.py)
; ============================================================
[program:dashboard]
command=python3 /app/installer/app.py
directory=/app
autostart=true
autorestart=true
priority=30
stopasgroup=true
stdout_logfile=/data/logs/dashboard.log
stderr_logfile=/data/logs/dashboard.err
environment=PYTHONUNBUFFERED="1",DISPLAY=":99",HOME="/root"

; ============================================================
; OpenClaw Gateway - Main multi-agent gateway
; ============================================================
[program:openclaw]
command=openclaw gateway --allow-unconfigured --port 3000 --bind lan --state-dir /data/openclaw
directory=/data/openclaw
autostart=true
autorestart=true
priority=40
stopasgroup=true
killasgroup=true
stdout_logfile=/data/logs/openclaw.log
stderr_logfile=/data/logs/openclaw.err
environment=NODE_ENV="production",OPENCLAW_STATE_DIR="/data/openclaw",NODE_OPTIONS="--max-old-space-size=4096",DISPLAY=":99",HOME="/root"
